# docker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
silent=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ;

[program:tianyianquan_script]
command=bash /tianyianquan/scripts/60-seconds-mysql-search.sh
directory=/tianyianquan
autostart=true
autorestart=true
redirect_stderr = true
stdout_logfile=/var/log/tianyianquan/tianyianquan_script.out.log
stdout_logfile_maxbytes = 20MB
# 日志文件最多备份数量
stdout_logfile_backups = 20

[program:tianyianquan_web]
command=java -jar /tianyianquan/app.jar --server.port=8080
directory=/tianyianquan
autostart=true
autorestart=true
redirect_stderr = true
stdout_logfile=/var/log/tianyianquan/tianyianquan_web.out.log
stdout_logfile_maxbytes = 20MB
# 日志文件最多备份数量
stdout_logfile_backups = 20